# dnsmasq.service
[Unit]
Description=dnsmasq - lightweight DHCP and caching DNS server
After=network.target
ConditionFileNotEmpty=/etc/dnsmasq.conf

[Service]
Type=simple
ExecStart=/usr/sbin/dnsmasq -k
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
