# wlan1-setup.service
[Unit]
Description=Setup wlan1 for AP mode (temperature_server)
After=network.target NetworkManager.service
Before=hostapd.service temperature-server.service

[Service]
Type=oneshot
# 環境変数（デフォルト値: 192.168.4.1/24, wlan1）
Environment=AP_IP=192.168.4.1
Environment=AP_SUBNET=24
Environment=AP_INTERFACE=wlan1
# setup_wlan1.shスクリプトを使用（環境変数を読み込む）
ExecStart=/usr/local/bin/setup_wlan1.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

